# spider.cron
*/2 * * * * python3 /noval/spider.py > /noval/noval.log
